# systemd/.config/systemd/user/amphora-sync.service
[Unit]
Description=Amphora vault git sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/amphora
ExecStart=/bin/bash -c '\
  git add -A && \
  git diff --cached --quiet || \
  git commit -m "vault backup: $(date +%%Y-%%m-%%d\ %%H:%%M:%%S)" && \
  git push origin master 2>/dev/null || git push origin main'
