# systemd/.config/systemd/user/omarchy-theme-auto.service
[Unit]
Description=Omarchy auto theme switch based on time of day

[Service]
Type=oneshot
ExecStart=%h/.local/bin/omarchy-theme-auto
